# scripts/backup.service
[Unit]
Description=Charger API Backup Service
After=network.target

[Service]
Type=oneshot
User=basar
Group=basar
WorkingDirectory=/home/basar/charger
Environment="PATH=/home/basar/charger/env/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/basar/charger/env/bin/python3 -u /home/basar/charger/scripts/backup_manager.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

